# dora.service
[Unit]
Description=Dora DHCP Server
Documentation=https://github.com/bluecatengineering/dora
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dora
Group=dora
WorkingDirectory=/var/lib/dora
ExecStart=/usr/bin/dora -c /etc/dora/config.yaml -d /var/lib/dora/leases.db
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/dora
ReadOnlyPaths=/etc/dora

# Network capabilities
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dora

[Install]
WantedBy=multi-user.target
